# luci-theme-material-old/root/etc/uci-defaults/30_luci-theme-material-old
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Material_Old=/luci-static/material-old
	set luci.main.mediaurlbase=/luci-static/material-old
	commit luci
EOF
exit 0
