# pbr/files/etc/uci-defaults/90-pbr
#!/bin/sh

uci -q batch <<-EOT
	delete firewall.pbr
	set firewall.pbr=include
	set firewall.pbr.type=script
	set firewall.pbr.path=/usr/share/pbr/pbr.firewall.include
	commit firewall
EOT

exit 0
